# 3rdparty/sgx-startup-ubuntu.conf
start on runlevel 2

script
	if [ -d /sys/devices/ocp.*/56000000.sgx ] ; then
		echo "sgx: Starting PVR"
	fi
end script
